# deploy/templates/default/uwsgi.conf.erb
# Emperor uWSGI Upstart script
# This file is managed by Chef for <%= node['fqdn'] %>

description "uWSGI"

start on [2345]
stop on runlevel [06]

script
    export HOME=/home/<%= @uid %>
    export NEW_RELIC_CONFIG_FILE=<%= @home_path %>/newrelic.ini
    <%= @home_path %>/bin/newrelic-admin run-program <%= @uwsgi_path %> -H <%= @home_path %> -s /tmp/uwsgi.<%= @name %>.sock --single-interpreter --env LC_ALL="en_US.UTF-8" --env LANG="en_US.UTF-8" --uid <%= @uid %> --gid <%= @gid %> --pidfile <%= @pid_path %> --die-on-term --processes=<%= @worker_processes %><%= @extra_params %> 2>&1
end script
