# deploy/templates/default/uwsgi.conf.erb
# Emperor uWSGI Upstart script
# This file is managed by Chef for <%= node['fqdn'] %>

description "uWSGI"

start on [2345]
stop on runlevel [06]

respawn
respawn limit 10 5

script
    UWSGI=<%= @uwsgi_path %>
    ROOT=<%= @home_path %>
    PID=<%= @pid_path %>

    cd $ROOT
    $UWSGI \
    --env LC_ALL="en_US.UTF-8" \
    --env LANG="en_US.UTF-8" \
    --uid <%= @uid %> \
    --gid <%= @gid %> \
    --pidfile $PID \
    --processes=<%= @worker_processes %><%= @extra_params %>
end script
